# sys-cluster/hadoop-common-bin/files/hadoop-namenode
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

# Source hadoop settings
. /etc/conf.d/hadoop

export HADOOP_CONF_DIR HADOOP_LOG_DIR YARN_LOG_DIR JAVA_HOME

# Operational commands
extra_commands="format"

depend() {
	use net
}

start() {
	ebegin "Starting Hadoop Namenode"
	start-stop-daemon --start --background --wait 100 --pidfile /run/${SVCNAME}.pid --make-pidfile --user hdfs --group hadoop --exec /usr/bin/hdfs -- namenode

	eend $?
}

stop() {
	ebegin "Stopping Hadoop Namenode"
	start-stop-daemon --stop --pidfile /run/${SVCNAME}.pid
	eend $?
}

format() {
	ebegin "Formatting HDFS"
	su hdfs -c "hdfs namenode -format"
	eend $?
}
