# sys-cluster/hadoop-common-bin/files/hadoop-datanode
#!/sbin/runscript
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

# Source hadoop settings
. /etc/conf.d/hadoop

export HADOOP_CONF_DIR HADOOP_LOG_DIR YARN_LOG_DIR JAVA_HOME

depend() {
	use net
}

start() {
	ebegin "Starting Hadoop Datanode"
	start-stop-daemon --start --background --wait 100 --pidfile /run/${SVCNAME}.pid --make-pidfile --user hdfs --group hadoop --exec /usr/bin/hdfs -- datanode

	eend $?
}

stop() {
	ebegin "Stopping Hadoop Datanode"
	start-stop-daemon --stop --pidfile /run/${SVCNAME}.pid
	eend $?
}
